# deploy/systemd/perfmatters-api.service
[Unit]
Description=Perfmatters Configuration Generator API
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/perfmatters-api
Environment=PATH=/opt/perfmatters-api/venv/bin
Environment=HOST=127.0.0.1
Environment=PORT=8080
Environment=DEBUG=False
Environment=FLASK_ENV=production
ExecStart=/opt/perfmatters-api/venv/bin/gunicorn -c gunicorn.conf.py app:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/perfmatters-api/logs
ReadWritePaths=/opt/perfmatters-api/generated_configs
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target